# web/supervisor/nodejs.conf
[program:nodejs]
command=node /docker/app.js
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0